# roles/manage/files/etc/profile.d/ruby.sh
export PATH=/usr/local/ruby-2.1.1/bin:$PATH
